# deployment/podman/404seo.service
# =============================================================================
# 404seo — Service systemd --user pour le boot auto (Podman 4.3.1, sans Quadlet)
# =============================================================================
# Quadlet n'existe pas en Podman 4.3.1 : on utilise un service oneshot qui lance
# start-containers.sh au demarrage. Le lingering (deja active pour seo) permet a
# ce service de tourner sans session SSH ouverte.
#
# Installation (en tant que seo) :
#   mkdir -p ~/.config/systemd/user
#   cp ~/404seo/deployment/podman/404seo.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now 404seo.service
#
# Le --restart unless-stopped des conteneurs gere les crashs ; ce service gere
# le demarrage au boot.
# =============================================================================

[Unit]
Description=404seo — stack Podman (site, api, redis, workers)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Demarre les 5 conteneurs
ExecStart=/usr/bin/bash %h/404seo/deployment/podman/start-containers.sh
# Arret propre : stoppe les conteneurs (sans les supprimer)
ExecStop=/usr/bin/podman stop 404seo-scheduler 404seo-worker-crawl 404seo-site 404seo-api 404seo-redis
TimeoutStartSec=300

[Install]
WantedBy=default.target
