# deployment/podman/quadlet/404seo-redis.container
# =============================================================================
# 404seo — Redis (Quadlet, rootless systemd --user)
# =============================================================================
# Installe via setup-podman.sh dans ~/.config/containers/systemd/
# Puis : systemctl --user daemon-reload && systemctl --user start 404seo-redis
# =============================================================================

[Unit]
Description=404seo Redis (BullMQ + cache)
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=404seo-redis
Image=docker.io/library/redis:7-alpine
Network=404seo.network
# Alias 'redis' : les conteneurs app le joignent via REDIS_HOST=redis
NetworkAlias=redis
Volume=404seo-redis-data.volume:/data
EnvironmentFile=%h/404seo/deployment/podman/.env
Exec=/bin/sh -c 'exec redis-server --requirepass "$REDIS_PASSWORD" --appendonly yes --maxmemory 256mb --maxmemory-policy allkeys-lru'
HealthCmd=/bin/sh -c 'redis-cli -a "$REDIS_PASSWORD" ping'
HealthInterval=10s
HealthRetries=5

[Service]
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
