# deployment/podman/quadlet/404seo-site.container
# =============================================================================
# 404seo — Next.js 16 site (Quadlet, rootless systemd --user, --network=host)
# =============================================================================
# En mode host, le site ecoute sur l'hote. BIND sur 127.0.0.1 (HOSTNAME) :
# seul le Nginx HestiaCP de l'hote y accede.
# =============================================================================
[Unit]
Description=404seo Next.js site
After=network-online.target 404seo-api.service
Wants=network-online.target
Requires=404seo-api.service

[Container]
ContainerName=404seo-site
Image=localhost/404seo-site:latest
Network=host
EnvironmentFile=%h/404seo/deployment/podman/.env
Environment=NODE_ENV=production
# bind loopback : seul Nginx (hote) accede au site
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3030
Exec=node site/server.js

[Service]
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
