# deployment/podman/quadlet/404seo-scheduler.container
# =============================================================================
# 404seo — BullMQ scheduler worker (Quadlet, rootless systemd --user, host)
# =============================================================================
# Reutilise l'image du worker (Playwright). Joint Redis/Postgres en 127.0.0.1.
[Unit]
Description=404seo BullMQ scheduler worker
After=network-online.target 404seo-redis.service 404seo-worker-crawl.service
Wants=network-online.target
Requires=404seo-redis.service

[Container]
ContainerName=404seo-scheduler
Image=localhost/404seo-worker:latest
Network=host
EnvironmentFile=%h/404seo/deployment/podman/.env
Environment=NODE_ENV=production
Exec=node_modules/.bin/tsx workers/schedule-worker.ts

[Service]
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
