# deployment/podman/quadlet/404seo-redis.container
# =============================================================================
# 404seo — Redis (Quadlet, rootless systemd --user, --network=host)
# =============================================================================
# En mode host, Redis ecoute directement sur l'hote. On le BIND sur 127.0.0.1
# pour qu'il ne soit JAMAIS accessible depuis l'exterieur (seuls les autres
# conteneurs host + l'hote y accedent).
# =============================================================================

[Unit]
Description=404seo Redis (BullMQ + cache)
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=404seo-redis
Image=docker.io/library/redis:7-alpine
Network=host
Volume=404seo-redis-data.volume:/data
EnvironmentFile=%h/404seo/deployment/podman/.env
# bind 127.0.0.1 : non expose publiquement. requirepass via $REDIS_PASSWORD.
Exec=/bin/sh -c 'exec redis-server --bind 127.0.0.1 --port 6379 --requirepass "$REDIS_PASSWORD" --appendonly yes --maxmemory 256mb --maxmemory-policy allkeys-lru'
HealthCmd=/bin/sh -c 'redis-cli -a "$REDIS_PASSWORD" ping'
HealthInterval=10s
HealthRetries=5

[Service]
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
