# deployment/podman/quadlet/404seo-worker-crawl.container
# =============================================================================
# 404seo — BullMQ crawl worker (Quadlet, rootless systemd --user)
# =============================================================================
# Image basee sur Playwright (Chromium + deps systeme inclus).
[Unit]
Description=404seo BullMQ crawl worker
After=network-online.target 404seo-redis.service
Wants=network-online.target
Requires=404seo-redis.service

[Container]
ContainerName=404seo-worker-crawl
Image=localhost/404seo-worker:latest
Network=404seo.network
EnvironmentFile=%h/404seo/deployment/podman/.env
Environment=NODE_ENV=production
# Chromium a besoin de memoire partagee
ShmSize=1g
# Resolution du Postgres natif de l'hote
AddHost=host.containers.internal:host-gateway
Exec=node_modules/.bin/tsx workers/crawl-worker.ts

[Service]
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
